# core/modules/network-online-slx-targets/data/etc/systemd/system/network-dns.target
[Unit]
Description=DNS Server available
