# core/modules/network-online-slx-targets/data/etc/systemd/system/network-gateway.target
[Unit]
Description=Default route configured
