# core/modules/safe-mode/data/etc/systemd/system/exam-mode.target
[Unit]
Description=bwLehrpool Klausurmodus
Requires=multi-user.target
After=multi-user.target
