# core/modules/safe-mode/data/etc/systemd/system/exam-mode.target
[Unit]
Description=bwLehrpool Prüfungsmodus
Requires=multi-user.target
After=multi-user.target
