# core/modules/safe-mode/data/etc/systemd/system/safe-mode.service
[Unit]
Description=Sets up safe mode (no magic keys, no tty switch, no openbox menu)
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-safe_mode
