# core/modules/beamergui/data/etc/udev/rules.d/99-beamergui.rules
# Tell beamergui to show up
#
SUBSYSTEM=="drm", ACTION=="change", RUN+="/opt/openslx/bin/beamergui -w"
