# core/modules/debug/data/etc/systemd/system/debug-udevd_logger.service
[Unit]
Description=Logger for udevd-messages via udevadmin monitor
DefaultDependencies=no
IgnoreOnIsolate=yes

[Service]
ExecStart=/opt/openslx/scripts/systemd-debug_udevd_logger

