# core/modules/lightdm-greeter-bwlp/data/etc/systemd/system/lightdm.service.d/download-news.conf
[Service]
ExecStartPre=-/opt/openslx/scripts/lightdm-download_news
