# core/modules/dhcpc-busybox/data/etc/systemd/system/lightdm.service.d/lock-hostname.conf
[Service]
ExecStartPre=-/opt/openslx/scripts/udhcpc-openslx --lock-hostname-updates
