# core/modules/consolekit/data/etc/systemd/system/console-kit-daemon.service
[Unit]
Description=ConsoleKit

[Service]
Type=dbus
BusName=org.freedesktop.ConsoleKit
ExecStart=/usr/sbin/console-kit-daemon --no-daemon

