# core/modules/kexec-reboot/data/etc/systemd/system/enable-kexec-reboot.service
[Unit]
Description=Enable kexec-reboot via ctrl-alt-del
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/sh -c "systemctl enable kexec.target"
