# core/modules/cpugovernor/data/etc/systemd/system/cpu-governor.service
[Unit]
Description=Sets the CPU governor

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-cpu_governor
