# core/modules/iptables-helper/data/etc/systemd/system/openslx-iptables.service
[Unit]
Description=OpenSLX iptables helper

[Service]
ExecStart=/opt/openslx/iptables/iptables-reloader

