# core/rootfs/rootfs-stage32/data/etc/systemd/system/network-gateway.target
[Unit]
Description=Default route configured
ConditionPathExists=/run/network/gateway-ready

