# core/rootfs/rootfs-stage32/data/etc/systemd/system/network-dns.target
[Unit]
Description=DNS Server available
ConditionPathExists=/run/network/dns-ready

