# core/rootfs/rootfs-stage32/data/etc/systemd/system/stage4.target
[Unit]
Requires=mount-stage4.service
After=mount-stage4.service
Description=OpenSLX Stage 4
RefuseManualStart=yes
