# core/rootfs/rootfs-stage32/data/etc/systemd/system/tmp.target
[Unit]
Description=/tmp has been mounted
Wants=setup-partitions.service
After=setup-partitions.service

