# core/rootfs/rootfs-stage32/data/etc/systemd/system/update-etc_issue.service
[Unit]
Description=Update /etc/issue
Before=getty.target shutdown.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/openslx-create_issue
RemainAfterExit=no
