# core/modules/polkit/templates/systemd.service
[Unit]
Description=the cancer that is called PolKit

[Service]
Type=dbus
BusName=org.freedesktop.PolicyKit1
ExecStart=%POLKITD% --no-debug

