# core/modules/qemu/data/etc/systemd/system/qemu.service
[Unit]
Description=Sets up the QEMU hypervisor
Requires=run-virt-env.service
After=run-virt-env.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/openslx/scripts/systemd-qemu_env start
