# core/modules/dhcpc-busybox/data/etc/systemd/system/udhcpc@.service
[Unit]
Description=DHCP Client for %i
After=systemd-tmpfiles-setup.service

[Service]
Type=simple
ExecStart=/opt/openslx/scripts/systemd-udhcpc %i
Restart=on-failure

