# core/modules/dhcpc-busybox/data/etc/systemd/system/udhcpc@.service
[Unit]
Description=DHCP Client for %i
After=systemd-tmpfiles-setup.service

[Service]
Type=forking
PIDFile=/run/udhcpc/udhcpc.%i.pid
ExecStart=/opt/openslx/scripts/systemd-udhcpc++ %i
ExecStopPost=/opt/openslx/bin/rm /run/udhcpc/udhcpc.%i.pid

