# core/modules/run-virt/data/etc/systemd/system/run-virt-env.service
[Unit]
Description=Sets up the virtualization environment
Wants=graphical.target
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/openslx/scripts/systemd-run_virt_env
