# core/modules/idleaction/data/etc/cron.d/openslx-idleaction
# Trigger scripts that checks idle status of machine and triggers actions

SHELL=/bin/ash
PATH=/usr/sbin:/usr/bin:/sbin:/bin:/opt/openslx/sbin:/opt/openslx/bin

*/2 *    * * *   root	/opt/openslx/scripts/idleaction-cron_script

