# core/modules/disk-partitions/data/etc/systemd/system/persistent.target
[Unit]
Description=Persistent partition has been mounted (if found)
