# core/modules/disk-partitions/data/etc/systemd/system/tmp.target
[Unit]
Description=/tmp is properly backed by HDD if available
