# core/modules/idleaction/data/etc/systemd/system/idleaction-set_suspend.service
[Unit]
Description=Reset idletime and setup next rtcwake before suspend
DefaultDependencies=no
RefuseManualStart=true
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-idleaction_init suspend
RemainAfterExit=no

