# core/modules/idleaction/data/etc/systemd/system/idle-daemon.service
[Unit]
Description=OpenSLX idle-daemon

[Service]
ExecStart=/opt/openslx/scripts/idle-daemon-launcher
Restart=on-failure
RestartSec=1
