# core/modules/idleaction/data/etc/systemd/system/idleaction-set_normal.service
[Unit]
Description=Setup next rtcwake
DefaultDependencies=no
RefuseManualStart=true

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-idleaction_init
RemainAfterExit=no

