# core/modules/pam-slx-plug/data/etc/systemd/system/user@.service.d/cleanup-umount.conf
[Service]
ExecStopPost=!/opt/openslx/scripts/systemd-session_cleanup %I
