# core/modules/run-virt/data/etc/systemd/system/system-check.service.d/50-run-virt-deps.conf
[Unit]
# our system-check hook needs stuff from these
After=run-virt-env.service mount-vm-store.service
