# core/modules/docker-ce/data/etc/systemd/system/docker.service.d/no-trivial-root.conf
[Unit]
Requires=docker-no-trivial-root.service
After=docker-no-trivial-root.service
