# core/modules/dnbd3-proxy-mode/data/etc/systemd/system/dnbd3-proxy.target
[Unit]
Description=DNBD3 server running
