# core/modules/dnbd3-proxy-mode/data/etc/systemd/system/dnbd3-logwatch.service
[Unit]
Description=DNBD3 monitor reporting to Satellite Server

[Service]
ExecStart=/opt/openslx/scripts/systemd-dnbd3_logwatch
Restart=always
RestartSec=5
TimeoutStopSec=2
