# core/modules/dnbd3-proxy-mode/data/etc/systemd/system/dnbd3-tune-network.service
[Unit]
Description=Tweak network stack of %i for maximum throughput
Wants=dhcpc@%i.service
After=dhcpc@%i.service network-interface@%i.service

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-tweak_network %I

