# core/modules/smartctl/data/etc/systemd/system/smartctl.service
[Unit]
Description=Check HDD health status

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/openslx/scripts/systemd-smartctl

