# core/modules/etherwake/data/etc/cron.d/enable-wol
SHELL=/bin/ash
PATH=/usr/sbin:/usr/bin:/sbin:/bin:/opt/openslx/sbin:/opt/openslx/bin

11 * 	* * * 	root 	/opt/openslx/scripts/cron-enable_wol
@reboot 			root 	/opt/openslx/scripts/cron-enable_wol
