# core/rootfs/rootfs-stage32/data/etc/systemd/system/setup-partitions.service.d/wait-udev.conf
[Unit]
After=systemd-udev-settle.service systemd-udevd.service
