# core/modules/systemd/data/usr/lib/systemd/system/remote-fs.target
[Unit]
Description=Remote File Systems Impostor

