# core/modules/consolekit/templates/consolekit.systemd.service
[Unit]
Description=ConsoleKit (Legacy)
After=syslog.target

[Service]
Type=dbus
BusName=org.freedesktop.ConsoleKit
ExecStart=%DAEMON% --no-daemon

