# core/modules/slx-brightness/data/etc/systemd/system/slx-brightness.service
[Unit]
Description=Handle display brightness keys

[Service]
Type=forking
ExecStart=/opt/openslx/scripts/systemd-slx_brightness
Restart=on-failure
