# core/modules/hardware-stats/data/etc/systemd/system/hardware-stats.service
[Unit]
Description=Gather hardware information about this machine
After=tmp.target
Wants=tmp.target system-check.service
Before=system-check.service

[Service]
Type=oneshot
ExecStart=/opt/openslx/scripts/systemd-hardware_stats
RemainAfterExit=yes

