# core/modules/kexec-reboot/data/etc/systemd/system/systemd-kexec.service.d/00-splash.conf
[Service]
ExecStartPre=-/opt/openslx/scripts/systemd-kexec_pre
